# modules/services/paperless/paperless-broker.container
[Unit]
Description=Paperless Redis Broker
After=network-online.target

[Container]
Image=docker.io/library/redis:8
ContainerName=paperless-broker
AutoUpdate=registry
PublishPort=127.0.0.1:63791:6379
Volume=/opt/services/paperless/redis:/data:Z

[Service]
Restart=always
User=paperless
Group=paperless

[Install]
WantedBy=multi-user.target